# callfwd.socket
[Unit]
Description=Call forwarder control socket

[Socket]
ListenDatagram=%t/callfwd
SocketMode=0600
RemoveOnStop=true
