# conf/callfwd-update-acl.timer
[Unit]
Description=Update callfwd ACLs

[Timer]
OnCalendar=minutely
Persistent=true
