# conf/callfwd-update-acl.service
[Unit]
Description=Update callfwd ACLs
Requires=callfwd.service

[Service]
Type=oneshot
ExecStart=/home/sergei/update_acl.sh
User=sergei
