# conf/callfwd-reload.timer
[Unit]
Description=CallFwd DB autoupdate
PartOf=callfwd.service

[Timer]
OnCalendar=daily
Persistent=true
