# conf/callfwd-acl.timer
[Unit]
Description=CallFwd ACLs autoupdate
PartOf=callfwd.service

[Timer]
OnCalendar=minutely
Persistent=true
