# conf/callfwd.socket
[Unit]
Description=CallFwd control socket

[Socket]
ListenDatagram=/run/callfwd
SocketMode=0666
RemoveOnStop=true
